# packaging/com.embitious.agent-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.embitious.agent-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.agm/Agent Monitor.app/Contents/MacOS/agm</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>AGM_PREFIX</key>
        <string>__HOME__/.agm</string>
        <key>NOTIFY_MAILBOX_PORT</key>
        <string>3847</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/AgentMonitor/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/AgentMonitor/stderr.log</string>
</dict>
</plist>
